# systemd/speedtest.crontab
*/5 * * * * /etc/zabbix_agentd.conf.d/bin/speedtest.sh
