# systemd/zabbix-speedtest-debian.service
# Autogenerated by generate-debian-files.sh. DO NOT EDIT MANUALLY. DO NOT INCLUDE IN PR.
[Unit]
Description=Run a speedtest
After=network.target

[Service]
Type=simple
User=zabbix
ExecStart=/etc/zabbix/bin/zbx-speedtest.sh --run

[Install]
WantedBy=multi-user.target
